# slim/packages/slim/postinst
#!/bin/bash

systemctl -f enable slim
systemctl set-default graphical.target
