# xdm/packages/xdm/postinst
#!/bin/bash

systemctl -f enable xdm
systemctl set-default graphical.target
