# xdm/packages/xdm/postinst
#!/bin/bash

if [ ! -L /etc/systemd/system/display-manager.service ]; then
systemctl -f enable xdm
systemctl set-default graphical.target
fi
