# xdm/packages/xdm/prerm
#!/bin/bash

systemctl -s stop xdm
systemctl -f disable xdm
systemctl set-default multi-user.target
