# logrotate/packages/logrotate/postinst
#!/bin/sh

systemctl -f enable logrotate
