# kernel/packages/cpupower/files
	/etc/cpufreq*.conf
	/etc/default/cpupower
	/usr/bin/cpufreq*
	/usr/bin/cpupower
	/usr/lib/libcpupower.so.*
	/usr/lib/systemd/scripts/cpupower
	/usr/lib/systemd/system/cpupower.service
	/usr/sbin/cpufreq*
	/usr/share/bash-completion
	/usr/share/locale/*/*/cpupower.mo
	/usr/share/man/man1/cpupower*.1
