# networkmanager/packages/networkmanager/postinst
#!/bin/bash

systemctl enable NetworkManager
