# gdm/packages/gdm/postinst
#!/bin/bash

#systemctl -f enable gdm
usermod -G gdm,video gdm
