# gdm/packages/gdm/prerm
#!/bin/bash

systemctl -s stop gdm
systemctl -f disable gdm
systemctl set-default multi-user.target
