# distcc/packages/distccd/files
	/etc/init.d/distccd
	/etc/runlevel.d/custom/distccd.service
	/etc/runlevel.d/default/distccd.service
	/etc/sysconfig/distccd
	/usr/bin/distccd
	/usr/lib/systemd/system/distccd.service
	/usr/share/man/man1/distccd.1
