# distcc/packages/distccd/postinst
#!/bin/bash

systemctl -f enable distccd
