# emptty/packages/emptty/postinst
#!/bin/bash

systemctl -f enable emptty
systemctl set-default graphical.target
