# modemmanager/packages/modemmanager/postinst
#!/bin/bash

systemctl enable ModemManager
