# postgresql/packages/postgresql-server/files
	/etc/init.d/postgresql
	/etc/postgresql
	/etc/runlevel.d/custom/postgresql.service
	/etc/runlevel.d/default/postgresql.service
	/usr/bin
	/usr/lib/libpgtypes.so
	/usr/lib/libpgtypes.so.*
	/usr/lib/postgresql
	/usr/lib/postgresql/bin
	/usr/lib/postgresql/dumpall
	/usr/lib/systemd/system
	/usr/share/doc
	/usr/share/man/man1
	/usr/share/man/man7
	/usr/share/postgresql/conversion_create.sql
	/usr/share/postgresql/extension/*.control
	/usr/share/postgresql/extension/*.sql
	/usr/share/postgresql/information_schema.sql
	/usr/share/postgresql/pg_hba.conf.sample
	/usr/share/postgresql/pg_ident.conf.sample
	/usr/share/postgresql/pg_service.conf.sample
	/usr/share/postgresql/postgres.bki
	/usr/share/postgresql/postgres.description
	/usr/share/postgresql/postgresql.conf.sample
	/usr/share/postgresql/postgres.shdescription
	/usr/share/postgresql/psqlrc.sample
	/usr/share/postgresql/recovery.conf.sample
	/usr/share/postgresql/snowball_create.sql
	/usr/share/postgresql/sql_features.txt
	/usr/share/postgresql/system_views.sql
	/usr/share/postgresql/timezone
	/usr/share/postgresql/timezonesets
	/usr/share/postgresql/tsearch_data
