# postgresql/packages/postgresql-server/prerm
#!/bin/bash

systemctl disable postgresql
