# bluez/packages/bluez/postinst
#!/bin/bash

systemctl enable bluetooth.service
