# lvm/packages/lvm/postinst
#!/bin/bash

systemctl enable blk-availability
systemctl enable lvm2-monitor
systemctl enable lvm2-lvmetad.socket
