# nvidia/addons/etc/xdg/autostart/nvidia-settings.desktop
[Desktop Entry]
Type=Application
Name=NVidia Settings
Exec=nvidia-settings -l
Terminal=false
X-KDE-StartupNotify=false
StartupNotify=false
