# rsync/packages/rsync/postinst
#!/bin/bash

systemctl -f enable rsyncd
