# audit/packages/audit/postinst
#!/bin/bash

if [ -f /usr/bin/systemctl ]; then
systemctl enable auditd
fi
