# sddm/packages/sddm/postinst
#!/bin/bash

systemctl -f enable sddm
systemctl set-default graphical.target
