# sddm/packages/sddm/prerm
#!/bin/bash

systemctl -s stop sddm
systemctl -f disable sddm
systemctl set-default multi-user.target
