# linuxconsole/packages/linuxconsole/postinst
#!/bin/bash

systemctl -f enable inputattach
