# linuxconsole/packages/linuxconsole/prerm
#!/bin/bash

systemctl disable inputattach
