# timidity/addons/usr/lib/systemd/system/timidity.service
[Unit]
Description=TiMidity++ Daemon
After=sound.target

[Service]
ExecStart=/usr/bin/timidity -iAD

[Install]
WantedBy=multi-user.target
