# dnsmasq/packages/dnsmasq/postinst
#!/bin/bash

systemctl -f enable dnsmasq
