# cups/packages/cups/postinst
#!/bin/bash

systemctl enable org.cups.cupsd
