# ipvsadm/packages/ipvsadm/postinst
#!/bin/bash

systemctl -f enable ipvsadm
