# lightdm/packages/lightdm/postinst
#!/bin/bash

if [ ! -L /etc/systemd/system/display-manager.service ]; then
systemctl -f enable lightdm
systemctl set-default graphical.target
fi
