# lightdm/packages/lightdm/postinst
#!/bin/bash

systemctl -f enable lightdm
systemctl set-default graphical.target
