# lightdm/packages/lightdm/prerm
#!/bin/bash

systemctl -s stop lightdm
systemctl -f disable lightdm
systemctl set-default multi-user.target
