# gvfs/packages/gvfs/postinst
#!/bin/bash

systemctl enable gvfs-daemon
