# avahi/packages/avahi/postinst
#!/bin/bash

systemctl enable avahi-daemon
