# subversion/packages/subversion/prerm
#!/bin/bash

systemctl disable svnserve
