# openldap/packages/openldap-server/files
	/etc/init.d
	/etc/openldap/DB_CONFIG.example
	/etc/openldap/schema
	/etc/openldap/slapd.ldif*
	/etc/runlevel.d
	/usr/lib/systemd/system/slapd.service
	/usr/sbin/slap*
	/usr/share/man/man5/slap*
	/usr/share/man/man8/slap*
	/var/lib/openldap
