# clamav/packages/clamav/postinst
#!/bin/bash

systemctl -f enable clamav-freshclam
systemctl -f enable clamav-daemon
