# open-vm-tools/packages/open-vm-tools/postinst
#!/bin/bash
systemctl -f enable vmtoolsd
