# pulseaudio/packages/pulseaudio/postinst
#!/bin/bash

systemctl --global enable pulseaudio.socket
