# pulseaudio/packages/pulseaudio/prerm
#!/bin/bash

systemctl --global disable pulseaudio.socket
