# lxdm/packages/lxdm/postinst
#!/bin/bash

systemctl -f enable lxdm
systemctl set-default graphical.target
