# lxdm/packages/lxdm/prerm
#!/bin/bash

systemctl -s stop lxdm
systemctl -f disable lxdm
systemctl set-default multi-user.target
