# cyrus-sasl2/addons/usr/lib/systemd/system/saslauthd.service
[Unit]
Description=SASL authentication daemon.

[Service]
Type=forking
ExecStart=/usr/sbin/saslauthd -a pam

[Install]
WantedBy=multi-user.target
