# cups-filters/packages/cups-filters/postinst
#!/bin/bash

systemctl -f enable cups-browsed
