# cups-filters/packages/cups-filters/prerm
#!/bin/bash

systemctl disable cups-browsed
