# packagekit/packages/packagekit/postinst
#!/bin/bash

systemctl -f enable packagekit
