# packagekit/packages/packagekit/prerm
#!/bin/bash

systemctl disable packagekit
