# optimus-manager/packages/optimus-manager/postinst
#!/bin/sh -eux

systemctl -f enable optimus-manager.service
