# gkrellm/packages/gkrellm/postinst
#!/bin/bash

systemctl -f enable gkrellm
