# wpasupplicant/packages/wpasupplicant/postinst
#!/bin/bash

systemctl -f enable wpa_supplicant
