# wpasupplicant/packages/wpasupplicant/postinst
#!/bin/bash

systemctl enable wpa_supplicant
