# wpasupplicant/packages/wpasupplicant/prerm
#!/bin/bash

systemctl disable wpa_supplicant
