# findutils/packages/locate/files
	/etc/cron.weekly/updatedb
	/etc/sysconfig/updatedb
	/usr/bin/locate
	/usr/lib/updatedb
	/usr/sbin/updatedb
	/usr/share/man/man1/locate.1
	/usr/share/man/man5/locatedb.5
	/usr/share/man/man8/updatedb.8
