# glibc/packages/nscd/files
	/etc/nscd.conf
	/usr/lib/systemd/system/nscd.service
	/usr/lib/tmpfiles.d/nscd.conf
	/usr/sbin/nscd
	/usr/sbin/nscd_nischeck
