# samba/packages/samba-client/files
	/etc/init.d/winbindd
	/etc/runlevel.d/custom/winbindd.service
	/etc/runlevel.d/default/winbindd.service
	/lib/libnss_winbind.so.2
	/lib/security/pam_smbpass.so
	/lib/security/pam_winbind.so
	/usr/bin/findsmb
	/usr/bin/net
	/usr/bin/nmblookup
	/usr/bin/ntlm_auth
	/usr/bin/profiles
	/usr/bin/rpcclient
	/usr/bin/smbcacls
	/usr/bin/smbclient
	/usr/bin/smbcquotas
	/usr/bin/smbget
	/usr/bin/smbtar
	/usr/bin/smbtree
	/usr/bin/wbinfo
	/usr/lib/smbclient.so.*
	/usr/sbin/mount.cifs
	/usr/sbin/umount.cifs
	/usr/sbin/winbindd
	/usr/share/man/man1/findsmb.1
	/usr/share/man/man1/nmblookup.1
	/usr/share/man/man1/ntlm_auth.1
	/usr/share/man/man1/profiles.1
	/usr/share/man/man1/rpcclient.1
	/usr/share/man/man1/smbcacls.1
	/usr/share/man/man1/smbclient.1
	/usr/share/man/man1/smbcquotas.1
	/usr/share/man/man1/smbget.1
	/usr/share/man/man1/smbtar.1
	/usr/share/man/man1/smbtree.1
	/usr/share/man/man1/wbinfo.1
	/usr/share/man/man5/smbgetrc.5
	/usr/share/man/man8/net.8
	/usr/share/man/man8/winbindd.8
