# nagios/packages/nagios/postinst
#!/bin/bash

systemctl -f enable nagios
