# cyrus-sasl2/packages/cyrus-sasl2/prerm
#!/bin/bash

systemctl disable saslauthd
