# accountsservice/packages/accountsservice/postinst
#!/bin/bash

systemctl -f enable accounts-daemon
