# accountsservice/packages/accountsservice/prerm
#!/bin/bash

systemctl disable accounts-daemon
