# kernel/packages/kernel/postinst
#!/bin/sh

systemctl enable icmpredirect_fix

exec mkinitrd-installedkernel
