# htpdate/packages/htpdate/postinst
#!/bin/bash

systemctl -f enable htpdate
